# sourcecode/web-ssh/webssh.ini
[program:webssh]
command = /home/prod/.virtualenvs/webssh/bin/python run.py --port=8567 --fbidhttp=False --logging=info --log-file-prefix=webssh.log

directory = /home/prod/deploys/webssh
autostart = true
autorestart = true
startsecs = 10

stdout_logfile = /home/prod/deploys/log/webssh.log
stderr_logfile = /home/prod/deploys/log/webssh.error.log
stdout_logfile_maxbytes = 64MB
stderr_logfile_maxbytes = 64MB
stdout_logfile_backups = 3
kill_as_group = true